# YianaOCRService/com.vitygas.yiana-ocr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vitygas.yiana-ocr</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/devon/bin/yiana-ocr</string>
        <string>watch</string>
        <string>--path</string>
        <string>/Users/devon/Library/Mobile Documents/iCloud~com~vitygas~Yiana/Documents</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/devon/Library/Logs/yiana-ocr.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/devon/Library/Logs/yiana-ocr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
